# roles/node_exporter/templates/node_exporter.service.j2
# =============================================================================
# Prometheus Node Exporter Systemd Service
# =============================================================================
# 由 Ansible 自动生成 | Generated by Ansible
# Node Exporter 系统指标收集服务 | Node Exporter system metrics collection service
# =============================================================================

[Unit]
Description=Prometheus Node Exporter
Description=Prometheus 节点导出器 - 收集系统指标
Documentation=https://github.com/prometheus/node_exporter
After=network-online.target

[Service]
Type=simple
User=node_exporter
Group=node_exporter

# Node Exporter 启动命令 | Node Exporter start command
# 监听端口 | Listen port: {{ prometheus.node_exporter.port }}
# 排除的文件系统挂载点 | Excluded filesystem mount points
# 忽略的网络设备 | Ignored network devices
ExecStart=/usr/local/bin/node_exporter \
    --web.listen-address=:{{ prometheus.node_exporter.port }} \
    --collector.filesystem.mount-points-exclude=^/(dev|proc|sys|var/lib/docker/.+|var/lib/kubelet/.+)($|/) \
    --collector.netclass.ignored-devices=^(veth.*)$ \
    --collector.netdev.device-exclude=^(veth.*)$

SyslogIdentifier=node_exporter
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
