# roles/anix_cloudflared/templates/cloudflared.service.j2
# =============================================================================
# Cloudflare Tunnel Systemd Service File
# =============================================================================
# 此服务由 Ansible Role 'anix_cloudflared' 自动生成
# This service is automatically generated by Ansible Role 'anix_cloudflared'
# =============================================================================

[Unit]
Description=Cloudflare Tunnel
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ anix_cloudflared_user }}
Group={{ anix_cloudflared_group }}

# 服务启动命令 | Service start command
ExecStart={{ anix_cloudflared_binary_path }} tunnel run \
  --token {{ cf_tunnel_token }} \
  {{ anix_cloudflared_tunnel_args | join(' ') }}

# 服务重启策略 | Service restart policy
Restart=on-failure
RestartSec=10s

# 安全加固 | Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ anix_cloudflared_log_dir }}

# 日志配置 | Logging configuration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflared

[Install]
WantedBy=multi-user.target
